# doormand.service
[Unit]
Description=doorman Face Authentication Daemon
After=network.target
Documentation=https://github.com/yourusername/doorman

[Service]
Type=simple
ExecStart=/usr/local/bin/doormand
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/doorman /run

# Environment
Environment="RUST_LOG=info"

[Install]
WantedBy=multi-user.target

